# packaging/deb/postinst
#!/bin/sh
set -e

# FilterREX Connector Host — Post-install
# Reloads systemd. Does NOT auto-enable or start the service —
# the admin must configure /etc/filterrex/connector.env first.

case "$1" in
  configure)
    # Reload systemd to pick up the new/updated unit file
    if [ -d /run/systemd/system ]; then
      systemctl daemon-reload >/dev/null 2>&1 || true
    fi

    # Ensure config file has correct permissions
    if [ -f /etc/filterrex/connector.env ]; then
      chown root:filterrex /etc/filterrex/connector.env
      chmod 0640 /etc/filterrex/connector.env
    fi

    echo ""
    echo "============================================================"
    echo "  FilterREX Connector Host installed successfully."
    echo ""
    echo "  Next steps:"
    echo "    1. Edit /etc/filterrex/connector.env"
    echo "       Set FILTERREX_ENROLLMENT_TOKEN (or CONNECTOR_TOKEN)"
    echo ""
    echo "    2. Enable and start the service:"
    echo "       sudo systemctl enable --now filterrex-connector"
    echo "============================================================"
    echo ""
    ;;
esac

#DEBHELPER#

exit 0
