# packaging/filterrex-connector.service
[Unit]
Description=FilterREX Connector Host
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=filterrex
EnvironmentFile=/etc/filterrex/connector.env
ExecStart=/usr/bin/filterrex-connector
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/filterrex
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
